# .config/autostart/keyboard.desktop
[Desktop Entry]
Terminal=false
Type=Application
StartupNotify=false
Name=Keyboard setting
Exec=/home/feng/keysetting.sh
